# scripts/systemd/voice-stt-ptt.service.in
[Unit]
Description=voice-stt push-to-talk evdev listener
Documentation=https://github.com/MaxInertia/claude-voice-input-channel
# Start after voice-sttd so the ctrl socket is ready when the listener
# tries to send its first command. The listener is still resilient if
# the daemon isn't up yet (commands fail with a stderr message and the
# listener keeps waiting), but ordering here avoids the noisy error.
After=voice-sttd.service
Requires=voice-sttd.service

[Service]
Type=simple
WorkingDirectory=@VOICE_STT_ROOT@
EnvironmentFile=-@VOICE_STT_ROOT@/.env
Environment=PATH=%h/.local/bin:%h/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=%h/.local/bin/uv run --project @VOICE_STT_ROOT@ voice-stt-ptt
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
